# agent-docker/hermes/supervisord.conf
[program:hermes]
command=/bin/bash -c "if [ -f /tmp/.observability_disabled ]; then exec env -u PYTHONPATH /opt/hermes/docker/entrypoint.sh gateway run; else exec /opt/hermes/docker/entrypoint.sh gateway run; fi"
user=root
environment=HOME="/root",PYTHONPATH="%(ENV_PYTHONPATH)s"
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startretries=-1
stopasgroup=true
killasgroup=true

[program:hermes-dashboard]
command=/bin/bash -c "while [ ! -f ${HERMES_HOME:-/opt/data}/gateway.pid ]; do sleep 2; done; sleep 3; if [ -f /tmp/.observability_disabled ]; then exec env -u PYTHONPATH /opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --insecure; else exec /opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --insecure; fi"
user=root
environment=HOME="/root",PYTHONPATH="%(ENV_PYTHONPATH)s"
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autorestart=true
stopasgroup=true
killasgroup=true
