# agent-docker/hermes/supervisord.conf
[program:hermes]
command=/opt/hermes/docker/entrypoint.sh gateway run
user=root
environment=HOME="/root"
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autorestart=true
startretries=-1
stopasgroup=true
killasgroup=true

[program:hermes-dashboard]
command=/bin/bash -c "while [ ! -f /root/.hermes/gateway.pid ]; do sleep 2; done; sleep 3; exec /opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --insecure"
user=root
environment=HOME="/root"
redirect_stderr=true
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
autorestart=true
stopasgroup=true
killasgroup=true
